# deploy/clow.service
[Unit]
Description=Clow AI Web App
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/clow
EnvironmentFile=/root/.clow/app/.env
Environment=PATH=/root/clow/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/root/clow/.venv/bin/uvicorn clow.webapp:app --host 0.0.0.0 --port 8001 --log-level info
Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clow

[Install]
WantedBy=multi-user.target
