# deploy/clow.service
[Unit]
Description=Clow AI Web App
After=network.target
Wants=network-online.target

[Service]
Type=exec
User=root
WorkingDirectory=/root/clow
Environment=PATH=/root/clow/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/root/clow/.venv/bin/uvicorn clow.webapp:app --host 0.0.0.0 --port 8001 --workers 2 --log-level info
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clow

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/root/clow /root/.clow /root/.claude /tmp
PrivateTmp=true

# Health check — restart if unhealthy
WatchdogSec=120

[Install]
WantedBy=multi-user.target
